# nasa_back.timer
[Unit]
Description=Daily run nasa_back

[Timer]
OnCalendar=*-*-* 00:00:00
Persistent=true

[Install]
WantedBy=timers.target
